# configs/etc/rc.local
#!/bin/sh

sh ~/bbb-ap/remote_scripts/post_boot.sh

exit 0
